# deploy/sandbox-runtime/supervisord.conf
[unix_http_server]
file=/run/cocola/supervisor.sock
chmod=0770
chown=root:cocola

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/run/cocola/supervisord.pid
childlogdir=/run/cocola

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/cocola/supervisor.sock

[program:code-server]
command=/opt/cocola/code-server-launch.sh
directory=/
autostart=true
# Startup failures remain bounded by startretries. Once the optional editor has
# reached RUNNING, keep a later crash visible as EXITED instead of entering an
# unbounded crash/restart loop that can consume the whole Sandbox.
autorestart=false
exitcodes=0
startsecs=2
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
